# build/sysroot/usr/lib/systemd/system/casaos-user-service.service
[Unit]
After=casaos-message-bus.service
Description=CasaOS User Service

[Service]
ExecStartPre=/usr/bin/casaos-user-service -v
ExecStart=/usr/bin/casaos-user-service -c /etc/casaos/user-service.conf
PIDFile=/var/run/casaos/user-service.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
